# package/etc/systemd/system/maufilelist.service
[Unit]
Description=A configurable file listing system.

[Service]
ExecStart=/usr/bin/maufilelist

[Install]
WantedBy=multi-user.target
